# systemd/djtoolkit.timer
[Unit]
Description=DJ Toolkit — run pipeline every 30 minutes
Documentation=https://github.com/koolkatdj/djtoolkit

[Timer]
# Run 5 minutes after boot (lets the desktop settle)
OnBootSec=5min
# Then run every 30 minutes
OnUnitActiveSec=30min
# Randomize start within ±2 minutes to avoid thundering herd
RandomizedDelaySec=2min
# Persist: if the machine was off during a scheduled run, run it on next boot
Persistent=true
Unit=djtoolkit.service

[Install]
WantedBy=timers.target
